# scripts/t-recur.service
[Unit]
Description=Todo Recur

[Service]
ExecStart=/home/jon/Documents/Settings/dotfiles/scripts/t-recur.sh

[Install]
WantedBy=multi-user.target
